# home/Library/LaunchAgents/com.mattcarr.chezmoi-autosync.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mattcarr.chezmoi-autosync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/opt/homebrew/bin/chezmoi re-add 2>&amp;1</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/chezmoi-autosync.log</string>

    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/chezmoi-autosync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>{{ .chezmoi.homeDir }}</string>
    </dict>
</dict>
</plist>
